# supervisor.conf
;===============================================
; NOADS - Supervisor Configuration
;===============================================
; Use este arquivo para gerenciar o serviço do NOADS
; 
; INSTALAÇÃO:
; 1. Instale supervisor: apt-get install supervisor
; 2. Copie este arquivo: cp supervisor.conf /etc/supervisor/conf.d/noads.conf
; 3. Recarregue: supervisorctl reread && supervisorctl update
; 4. Inicie: supervisorctl start noads
;
; COMANDOS ÚTEIS:
; - supervisorctl status           # Ver status
; - supervisorctl restart noads    # Reiniciar
; - supervisorctl stop noads       # Parar
; - tail -f /var/log/supervisor/noads.log

[program:noads]

; Comando para executar
command=/home/seu_usuario/noads/venv/bin/python3 /home/seu_usuario/noads/simple_api.py

; Diretório de trabalho
directory=/home/seu_usuario/noads

; Qual usuário vai rodar
user=seu_usuario

; Auto-iniciar
autostart=true

; Reiniciar se cair
autorestart=true

; Quantas vezes tentar reiniciar
startretries=3

; Logs
stdout_logfile=/home/seu_usuario/noads/logs/noads.log
stderr_logfile=/home/seu_usuario/noads/logs/noads_error.log

; Tamanho máximo de log (10MB)
stdout_logfile_maxbytes=10485760
stdout_logfile_backups=10
stderr_logfile_maxbytes=10485760
stderr_logfile_backups=10

; Timeout para parar o processo
stopwaitsecs=10

; Período de espera antes de reiniciar (se cair)
startsecs=10
